# contents/etc/xdg/autostart/nvidia-prime.desktop
[Desktop Entry]
Type=Application
Name=Nvidia PRIME offloading
Exec=/usr/bin/nvidia-offload
NoDisplay=true
X-GNOME-Autostart-Phase=DisplayServer
